# services/aqbota.service
[Unit]
Description=Aqbota_roslaunch

[Service]
User=root
ExecStart=/home/robot/aqbota.sh

[Install]
WantedBy=multi-user.target